# templates/default/systemd.service.erb
[Unit]
Description=<%= @name %>
After=<%= @options['after'].is_a?(Array) ? @options['after'].join(' ') : @options['after'] %>

[Service]
Type=simple
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
User=<%= @user %>
KillSignal=TERM
MemoryLimit=1G
Restart=always
RestartSec=20
StartLimitInterval=0
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
